# roles/router/files/router.service
[Unit]
Description=Router

[Service]
Type=oneshot
ExecStart=/sbin/router.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
